# docker/postgresql/backup/crontab
0 1 * * * cd ~/projects/briefmetrics/backup/auto && ../autobackup.sh
